# src/etc/skel/.config/autostart/lucidsystems-greeter.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/lucidsystems-greeter
Hidden=false
X-GNOME-Autostart-enabled=true
Name=LucidSystems Greeter
Comment=
